# hyper_ws/grid.service
[Unit]
Description=Hyperliquid spot grid engine
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
WorkingDirectory=/home/youruser/kn/hyper_ws
ExecStart=/usr/bin/env python3 /home/youruser/kn/hyper_ws/grid.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
